# etc/init.d/irq-pin
#!/sbin/openrc-run

description="Pin device IRQs to selected CPUs based on /etc/conf.d/irq-pin"

depend() {
    need localmount
    after modules
}

load_cfg() {
    [ -r /etc/conf.d/irq-pin ] && . /etc/conf.d/irq-pin
}

pin_irqs_for() {
    local dev="$1" cpulist="$2" node
    
    awk -v d="$dev" -F: '
        /^[ 0-9]+:/ {irq=$1} 
        {rhs=substr($0,index($0,":")+1)} 
        rhs ~ "(^|[^A-Za-z0-9_])" d "([^A-Za-z0-9_]|$)" {print irq}
    ' /proc/interrupts | while read -r node; do
        [ -n "$node" ] || continue
        
        echo "$cpulist" > "/proc/irq/$node/smp_affinity_list" 2>/dev/null || \
            ewarn "Failed to set affinity for IRQ $node ($dev)"
    done
}

start() {
    ebegin "Pinning device IRQs"
        
    einfo "Refreshing IRQ pin configuration..."
    if command -v autopin-irqs >/dev/null 2>&1; then
        autopin-irqs || ewarn "autopin-irqs failed; using existing configuration."
    else
        ewarn "autopin-irqs not found in PATH; cannot generate configuration."
    fi

    load_cfg
    
    if [ -z "${DEVICES:-}" ]; then
        ewarn "No DEVICES defined in /etc/conf.d/irq-pin. Nothing to do."
        eend 0
        return
    fi
    
    local oifs="$IFS"; IFS=','
    for d in $DEVICES; do
        local var="CPUSET_${d}" cpus
        eval cpus=\${$var}
        if [ -n "$cpus" ]; then
            veinfo "Pinning $d to CPUs: $cpus"
            pin_irqs_for "$d" "$cpus"
        else
            ewarn "No CPUSET_$d found for device $d"
        fi
    done
    IFS="$oifs"
    eend 0
}

stop() {
    ebegin "Resetting IRQ affinities to online CPUs"
    load_cfg
    [ -n "${DEVICES:-}" ] || { eend 0; return; }
    
    local all_cpus
    all_cpus="$(cat /sys/devices/system/cpu/online 2>/dev/null || echo 0-$(( $(nproc) - 1 )))"
    
    veinfo "Resetting IRQs for $DEVICES to CPUs ($all_cpus)"
    local oifs="$IFS"; IFS=','
    
    for d in $DEVICES; do
        awk -v dev="$d" -F: '
            /^[ 0-9]+:/ {irq=$1} 
            {rhs=substr($0,index($0,":")+1)} 
            rhs ~ "(^|[^A-Za-z0-9_])" dev "([^A-Za-z0-9_]|$)" {print irq}
        ' /proc/interrupts | while read -r node; do
            [ -n "$node" ] || continue
            
            # CORRECT: Write *only* to smp_affinity_list
            echo "$all_cpus" > "/proc/irq/$node/smp_affinity_list" 2>/dev/null
        done
    done
    
    IFS="$oifs"
    eend 0
}
