# etc/init.d/mountbinds
#!/sbin/openrc-run
# Copyright 2024 Uri Herrera <uri_herrera@nxos.org>
# Distributed under the terms of the BSD-3-Clause license

description="Creates FHS bind mounts for the Nitrux filesystem hierarchy"

depend() {
	# This script binds directories like /home, /etc, and /var.
	# It absolutely must run after all local filesystems are mounted.
	need localmount
}

# Helper to check if this is a Live session.
_is_live_session() {
	[ -f /usr/share/initramfs-tools/scripts/casper-bottom/12fstab ]
}

start() {
	if _is_live_session; then
		einfo "Live session detected, skipping FHS bind mounts."
		return 0
	fi

	ebegin "Binding FHS directories"

	mount --bind /dev /System/Devices
	mount --bind /boot /System/Boot
	mount --bind /etc /System/Configuration
	mount --bind /usr/sbin /System/Binaries/Administrative
	mount --bind /usr/bin /System/Binaries/Standard
	mount --bind /opt /System/Binaries/Optional
	mount --bind /mnt /System/Filesystems/Internal
	mount --bind /media /System/Filesystems/External
	mount --bind /usr/lib /System/Libraries
	mount --bind /usr/share /System/Shareable
	mount --bind /var /System/Variable
	mount --bind /home /Users

	eend $?
}

stop() {
	# Unmounting bind mounts is handled globally by the 'umountfs'
	# shutdown service. A manual stop action is not needed here.
	return 0
}
