# etc/init.d/rapl-power
#!/sbin/openrc-run

description="Sets read permissions for the Intel RAPL energy file"

depend() {
	need sysfs
}

start() {
	local rapl_file="/sys/class/powercap/intel-rapl:0/energy_uj"

	if [ ! -f "${rapl_file}" ]; then
		einfo "Intel RAPL interface not found, skipping."
		return 0
	fi

	ebegin "Setting permissions for Intel RAPL energy file"
	chmod o+r "${rapl_file}"
	eend $?
}

stop() {
	return 0
}
