# etc/init.d/scx-dynamic-scheduler
#!/sbin/openrc-run

description="Auto-selects scx_bpfland (Laptop) or scx_lavd (Desktop)"

supervisor=supervise-daemon
pidfile="/run/scx_dynamic.pid"

depend() {
    need modules
    need sysfs
}

start_pre() {
    if ls /sys/class/power_supply/BAT* >/dev/null 2>&1; then
        command="/usr/bin/scx_bpfland"
        command_args=""
        mode="Laptop"
    else
        command="/usr/bin/scx_lavd"
        command_args="--performance"
        mode="Desktop"
    fi

    if [ ! -x "$command" ]; then
        eerror "Scheduler binary not found: $command"
        return 1
    fi

    einfo "Selected SCX Scheduler (${mode} Mode): $(basename "$command")"
}
