# etc/init.d/thp-settings
#!/sbin/openrc-run
# Copyright 2024 Uri Herrera <uri_herrera@nxos.org>
# Distributed under the terms of the BSD-3-Clause license

description="Configures Transparent Huge Page settings"

depend() {
	# This script modifies files in /sys, so it needs the sysfs
	# service to have run first.
	need sysfs
}

start() {
	local thp_path="/sys/kernel/mm/transparent_hugepage"

	if [ ! -d "${thp_path}" ]; then
		einfo "Transparent Huge Page support not found, skipping."
		return 0
	fi

	ebegin "Applying Transparent Huge Page settings"
	echo "always" > "${thp_path}/enabled"
	echo "0" > "${thp_path}/khugepaged/defrag"
	eend $? "Failed to apply Transparent Huge Page settings"
}

stop() {
	# This is a one-shot service, there is nothing to stop.
	return 0
}
